# config/supervisord/conf.d/mysql-backup-log.conf
[program:mysql-backup-log]
command=/bin/bash -c "tail -q -n 0 -F /mysql-backup.log"
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
